# contrib/com.pulpo.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.pulpo.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/pulpod</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/pulpod.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/pulpod.stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>RUST_LOG</key>
        <string>pulpod=info</string>
    </dict>
</dict>
</plist>
